# ddns/crontab
*/1  *  *  *  * pi cd /home/pi && bash ddupdate.sh
